# assets/.config/autostart/Compton.desktop
[Desktop Entry]
Name=Compton
Comment=Window Compositor
Exec=compton --dbus
Type=Application
StartupNotify=false
Terminal=false
OnlyShowIn=Openbox;XFCE;