# assets/.config/autostart/pywal.desktop
[Desktop Entry]
Name=pywal restorer
Comment=Restore pywal theme
Exec=wal -R
Type=Application
StartupNotify=false
Terminal=false
