# home/.config/autostart/cornora.desktop
[Desktop Entry]
Name=Cornora
Comment=Hotcorner Daemon
Exec=cornora -tl skippy-xd
Type=Application
StartupNotify=false
Terminal=false
OnlyShowIn=Openbox;XFCE;