# home/.config/autostart/xrdb.desktop
[Desktop Entry]
Name=xrdb
Comment=Restore X resource database
Exec=xrdb -load .Xresources
Type=Application
StartupNotify=false
Terminal=false
